# system_files/usr/lib/systemd/user/dms-niri-config.service
[Unit]
Description=Adds the standard dms niri configutation files

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/dms-niri-config

